# .github/supervisor/youdaonotelm-server.conf
# Supervisor 配置文件示例
# 文件名: youdaonotelm-server.conf
# 安装位置: /etc/supervisor/conf.d/youdaonotelm-server.conf
# 作用: 管理 youdaonotelm-server 进程的启动、停止、重启和日志

[program:youdaonotelm-server]
# 程序启动命令
command=/home/flandern/youdaonotelm/youdaonotelm-server

# 程序工作目录
directory=/home/flandern/youdaonotelm

# 运行用户（根据安全要求调整，建议使用非 root 用户）
user=root

# 是否自动启动
autostart=true

# 是否自动重启
autorestart=true

# 启动后等待时间（秒），用于判断启动是否成功
startsecs=10

# 启动失败重试次数
startretries=3

# 是否将 stderr 重定向到 stdout
redirect_stderr=true

# 日志文件路径
stdout_logfile=/home/flandern/youdaonotelm/logs/server.log

# 单个日志文件最大大小（字节）
stdout_logfile_maxbytes=50MB

# 保留的日志文件数量
stdout_logfile_backups=10

# 环境变量配置
# 根据你的应用需求添加或修改
environment=GIN_MODE="release",GO_ENV="production"

# 进程优先级（数值越小优先级越高）
priority=999

# 是否在启动前终止之前的进程
killasgroup=true
stopasgroup=true

# 停止信号（SIGTERM 优雅停止，SIGKILL 强制停止）
stopsignal=SIGTERM

# 停止等待时间（秒）
stopwaitsecs=30
